# services/terrarium-monitor.target
[Unit]
Description=Terrarium Monitoring System
Requires=multi-user.target
After=multi-user.target

[Install]
WantedBy=multi-user.target
